# etc/systemd/system/gpioset.service
[Unit]
Description=GPIO initial setup

[Service]
Type=idle
ExecStart=/usr/bin/python /root/gpio.py set

[Install]
WantedBy=multi-user.target
