# etc/systemd/system/gpiooff.service
[Unit]
Description=Before shutdown, turn off GPIO

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/usr/bin/python /root/gpiooff.py

[Install]
WantedBy=multi-user.target
