# etc/systemd/system/gpioset.service
[Unit]
Description=GPIO initial setup and poweroff

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/python /root/gpio.py set
ExecStop=/usr/bin/python /root/gpiooff.py

[Install]
WantedBy=multi-user.target
